# ollama/LaunchDeamon/com.ollama.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.ollama.serve</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/ollama</string>
            <string>serve</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <true/>
        <key>StandardOutPath</key>
        <string>/Users/robertw/.ollama/logs/ollama.log</string>
        <key>StandardErrorPath</key>
        <string>/Users/robertw/.ollama/logs/ollama-error.log</string>
        <key>EnvironmentVariables</key>
        <dict>
            <!-- Core‑Params -->
            <key>HOME</key>
            <string>/Users/robertw</string>
            <key>OLLAMA_HOST</key>
            <string>0.0.0.0:11434</string>
            <key>OLLAMA_ORIGINS</key>
            <string>*</string>
            <key>OLLAMA_KEEP_ALIVE</key>
            <string>86400</string>
            <key>OLLAMA_DEVICE</key>
            <string>metal</string>
            <key>OLLAMA_MAX_MEMORY</key>
            <string>13GiB</string>
            <key>PATH</key>
            <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        </dict>
    </dict>
</plist>
